# etc/pika.conf
description "pika upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid jove
exec /bin/sh /home/jove/jove/launch.sh